# images/consul/v1/files/systemd/format-data-volume.service
[Unit]
Description=Format EBS volume if needed
# ConditionFirstBoot=yes

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/usr/local/bin/format_EBS.sh
TimeoutStartSec=60